# remote-check/remote-check.service.example
[Unit]
Description=monitoring remote servers for connectivity via remote-check.sh.

[Service]
Type=simple
ExecStart=/absolute/path/to/script/remote-check.sh
RestartSec=10m
