# remote-check/remote-check.service.example
[Unit]
Description=Monitors uptime via pinging the public IP.

[Service]
Type=simple
ExecStart=/absolute/path/to/script/remote-check.sh
RestartSec=10m
