# local-check/local-check.service.example
[Unit]
Description=Monitors uptime via pinging the public IP.

[Service]
Type=simple
ExecStart=/absolute/path/to/script/local-check.sh
RestartSec=10m
